# driver/linux/tjekind_service.conf
start on startup
task
exec /usr/bin/python /var/www/Tjekind/driver/linux/nfc-reader3.py
